# debian/rust-spooler.service
[Unit]
Description=FLUX ERP to CUPS Print Server Bridge
After=network.target cups.service
Wants=cups.service

[Service]
Type=simple
User=flux-spooler
Group=flux-spooler
ExecStart=/usr/bin/rust-spooler run
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/flux-spooler

[Install]
WantedBy=multi-user.target
